# deploy/kb-backend.service
[Unit]
Description=Knowledge Base Backend (FastAPI + Uvicorn)
After=network.target

[Service]
Type=simple
# 运行用户（生产环境推荐专用用户）
User=kbuser
Group=kbuser
WorkingDirectory=/home/kbuser/projects/knowledge-base/backend

# 启动命令（多 worker 模式）
ExecStart=/home/kbuser/projects/knowledge-base/venv/bin/python start_prod.py \
    --workers 4 \
    --host 127.0.0.1 \
    --port 8000 \
    --log-level info

# 进程优雅关闭（发送 SIGTERM 后等 30 秒）
KillSignal=SIGTERM
TimeoutStopSec=30

# 自动重启
Restart=always
RestartSec=5
StartLimitBurst=10
StartLimitIntervalSec=60

# 资源限制（性能优化）
LimitNOFILE=65535
LimitNPROC=65535
LimitAS=4G          # 地址空间限制
LimitRSS=8G         # 常驻内存限制

# 环境变量
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONIOENCODING=utf-8
Environment=ENV=production
Environment=ALLOWED_ORIGINS=https://your-domain.com

# 日志输出到 journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kb-backend

# 安全加固
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/kbuser/projects/knowledge-base/backend/data

[Install]
WantedBy=multi-user.target